# data/xdg-desktop-portal-xapp-cinnamon.service.in
[Unit]
Description=Portal service (XApp Cinnamon implementation)

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.cinnamon
ExecStart=@libexecdir@/xdg-desktop-portal-xapp-cinnamon
